# feeder/adsb-tunnel.service
[Unit]
Description=Reverse SSH tunnel to VPS (HTTP + SSH)
After=network-online.target
Wants=network-online.target

[Service]
User=ads-b
ExecStart=/usr/bin/autossh -M 0 -N \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o StrictHostKeyChecking=no \
  -o ExitOnForwardFailure=yes \
  -i /home/ads-b/.ssh/id_adsb_vps \
  -R 0.0.0.0:HTTP_PORT:localhost:80 \
  -R 0.0.0.0:SSH_PORT:localhost:22 \
  VPS_USER@VPS_IP
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
